# includes/etc/xdg/autostart/create-flatpak-aliases.desktop
[Desktop Entry]
Name=Creates Flatpak Aliases
Exec=/usr/bin/create-flatpak-aliases
StartupNotify=false
Type=Application
NoDisplay=true
