# includes/etc/xdg/autostart/op-browser-support.desktop
[Desktop Entry]
Name=1Password Browser-Support
Exec=/usr/bin/op-browser-support
StartupNotify=false
Type=Application
NoDisplay=true
